# host/etc/systemd/system/backup-watchdog.service
[Unit]
Description=Alert if restic backups have gone stale
Wants=network-online.target
After=network-online.target
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/backup-watchdog.sh
