# host/etc/systemd/system/gluetun-verify.service
[Unit]
Description=Verify VPN downloaders are not leaking
After=docker.service
Wants=network-online.target
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/gluetun-verify.sh
