# host/etc/systemd/system/media-archive.service
[Unit]
OnFailure=ntfy-failure@%n.service
Description=Tier2 media cold sync (whole library -> Glacier Deep Archive)
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/media-archive.sh
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
